# inkbox-webhook.service
[Unit]
Description=Inkbox Email Webhook Receiver
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/openclaw-config
ExecStart=/home/ec2-user/openclaw-config/.venv/bin/inkbox-webhook
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
