# inkbox-media-stream.service
[Unit]
Description=Inkbox Phone Media Stream WS Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/openclaw-config
ExecStart=/home/ec2-user/openclaw-config/.venv/bin/inkbox-phone-media
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
